# piku-nginx@.path
[Unit]
Description=Monitor .piku/nginx for changes (%i)

[Path]
PathChanged=/home/%i/.piku/nginx
Unit=piku-nginx.service

[Install]
WantedBy=multi-user.target
