# uwsgi-piku.service
[Unit]
Description=Piku uWSGI Emperor
After=syslog.target

[Service]
ExecStart=/usr/local/bin/uwsgi-piku --ini /home/piku/.piku/uwsgi/uwsgi.ini
User=piku
Group=www-data
RuntimeDirectory=uwsgi-piku
Restart=always
KillSignal=SIGQUIT
Type=notify
StandardError=syslog
NotifyAccess=all

[Install]
WantedBy=multi-user.target
